# crontab
* * * * * echo -e "$redПривет из крона !!! $nc" >/proc/1/fd/1 2>/proc/1/fd/2
* * * * * /scripts/task.sh
